# migration/systemd/cas-serve@.service
# migration/systemd/cas-serve@.service
#
# Versioned systemd unit template for the cas HTTP bridge server. Phase 2
# of the Petrastella → Hetzner migration (cas-b5f1) installs this file to
# /etc/systemd/system/cas-serve@.service on daniel@87.99.156.244 and enables
# cas-serve@daniel. Edit this file in cas-src; re-run
# migration/phase2-prepare-target.sh to push updates to the server.
#
# Spec deviation: the task spec originally called for
#   EnvironmentFile=/home/%i/.config/cas/env
# but on the target server `.config/cas/env` is already occupied by
# cas-fb43's user-level token vault (GH_TOKEN, VERCEL_TOKEN, NEON_API_KEY,
# etc. — login-shell CLI creds). Per supervisor's option-(b) resolution on
# cas-b5f1, this unit uses a separate file `serve.env` scoped strictly
# to the systemd service so principle of least privilege is honored
# (cas-serve never inherits user CLI tokens) and future vault edits don't
# implicitly mutate the service environment.
#
# ExecStart uses `cas bridge serve` (not `cas serve`, which is the
# stdio MCP subcommand). The bearer token is read from CAS_SERVE_TOKEN in
# serve.env via clap's env-var binding — the token never appears in argv,
# so /proc/<pid>/cmdline, `systemctl status`, and `ps -ef` no longer expose
# it. Requires cas-cli from commit adding env = "CAS_SERVE_TOKEN" to
# ServeArgs::token (cas-2aa6 fix, 2026-04-11). Older cas binaries required
# --token ${CAS_SERVE_TOKEN} on the ExecStart line; that is no longer safe.

[Unit]
Description=CAS HTTP bridge server for %i
After=network.target
AssertPathExists=/home/%i/.config/cas/serve.env

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
EnvironmentFile=/home/%i/.config/cas/serve.env
ExecStart=/usr/local/bin/cas bridge serve --bind 127.0.0.1 --port 18999
Restart=on-failure
RestartSec=5

# Discard stdout so the cas startup banner's 'Token: <value>' line does NOT
# land in the systemd journal. The banner is informational and the token is
# already in serve.env for anyone who needs it. Errors still flow through
# stderr → journal for real operational signal.
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
